# Public/Final/tor-0.3.3.9/DEBIAN/postinst
#!/bin/sh

launchctl load -w /Library/LaunchDaemons/com.mc.tor.plist

exit 0
