# Public/Final/Packages/postgresql-10.4/Library/LaunchDaemons/pgsql.stop.switch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>pgsql.stop</string>
	<key>LowPriorityIO</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/libexec/stop_pgsql.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>WatchPaths</key>
	<array>
		<string>/private/var/mobile/Library/Preferences/pgsql.plist</string>
	</array>
</dict>
</plist>
